# Guiões/S10/demon.service
[Unit]
Description=My Daemon

[Service]
ExecStart=/home/zezao/2324-G03/Guioes/S10/demon

[Install]
WantedBy=multi-user.target
